# Debugging/libbpf-tools/resmon/resmon-exporter.service.in
[Unit]
Description=mlxsw resource monitor exporter
Requisite=resmon.service
After=resmon.service

[Service]
EnvironmentFile=-@SYSCONFDIR@/sysconfig/resmon-exporter
ExecStart=@BINDIR@/resmon-exporter.py $PARAMS
